# net-misc/tcp-proxy/files/tcp_proxy.initd
#!/sbin/openrc-run
# Copyright 2016-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="TCP Proxy"
description="TCP Proxy"
pidfile=${pidfile:-"/run/${SVCNAME}.pid"}
command_user=${command_user:-nobody:dialout}

command="/usr/bin/tcp-proxy"
command_background=true
retry="${retry:-TERM/20/KILL/5}"

depend() {
	need net
	after net
}
