# app-admin/javaapp-manager/files/0.1.4/initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Java application: ${RC_SVCNAME}"

pidfile=${pidfile:-"/run/${RC_SVCNAME}.pid"}
command_user=${user:-${RC_SVCNAME}}

logfile=${logfile:-/var/log/${RC_SVCNAME}/start.log}
output_log=${output_log:-${logfile}}
error_log=${error_log:-${logfile}}

directory=${directory:-/opt/${RC_SVCNAME}}

java_opts_file=${java_opts_file:-${directory}/jvm.options}
if [ ! -z "${java_opts_file}" ] && [ -e "${java_opts_file}" ]; then
  for opt in `grep "^-" ${java_opts_file}`
  do
    java_opts="${java_opts} ${opt}"
  done
fi

free_memory=`free -m | awk '/Mem:/ {print $2}'`
max_memory=`expr $free_memory '*' 87 '/' 100`
cpu_cores=`egrep -c 'processor([[:space:]]+):.*' /proc/cpuinfo`
ng_memory=`expr $cpu_cores '*' 150`
ng_memory_max=`expr $max_memory '/' 10`
if [ "$ng_memory" -gt "$ng_memory_max" ]
then
  ng_memory=$ng_memory_max
fi

[[ "${java_mem_min}" == "auto" ]] && java_mem_min="${max_memory}m"
[[ "${java_mem_max}" == "auto" ]] && java_mem_max="${max_memory}m"
[[ "${java_mem_new}" == "auto" ]] && java_mem_new="${ng_memory}m"

java_opts="${java_opts} ${java_mem_max:+-Xmx}${java_mem_max} ${java_mem_min:+-Xms}${java_mem_min} ${java_mem_new:+-Xmn}${java_mem_new}"

directory=${directory:-/opt/}

command=${command:-java}
command_args=${command_args:-${java_opts} -jar ${jar:-${directory}/${RC_SVCNAME}.jar}}
command_background="true"

ulimit -l unlimited
ulimit -n ${fd_limit:-1000000}

depend() {
  after net
}

